# deploy_tools/gunicorn-systemd.template.socket
[Unit]
Description=SITENAME socket

[Socket]
ListenStream=/run/SITENAME/socket

[Install]
WantedBy=sockets.target